# rootfs/etc/s6-overlay/s6-rc.d/nordvpn-setup/up
#!/usr/bin/with-contenv bash
# s6 service: One-time Meshnet setup
set -euo pipefail
log() { echo "[meshnet-setup] $1"; }
log "Waiting for NordVPN daemon..."
for i in $(seq 1 30); do
    if nordvpn status >/dev/null 2>&1; then
        log "Daemon is ready"
        break
    fi
    sleep 1
done
TOKEN=$(bashio::config 'token' || echo "")
NICKNAME=$(bashio::config 'nickname' || echo "")
ALLOW_ROUTING=$(bashio::config 'allow_routing' || echo "*")
ALLOW_LOCAL=$(bashio::config 'allow_local_network' || echo "*")
ALLOW_FILESHARE=$(bashio::config 'allow_fileshare' || echo "*")
ALLOW_REMOTE=$(bashio::config 'allow_remote' || echo "*")
if [ -z "$TOKEN" ]; then
    log "ERROR: No access token configured."
    exit 1
fi
log "Authenticating..."
if ! nordvpn login --token "$TOKEN" 2>&1; then
    log "ERROR: Auth failed"
    exit 1
fi
log "Authentication successful"
log "Enabling Meshnet..."
if ! nordvpn set meshnet on 2>&1; then
    log "ERROR: Failed to enable Meshnet"
    exit 1
fi
log "Meshnet enabled"
if [ -n "$NICKNAME" ]; then
    log "Setting nickname: $NICKNAME"
    nordvpn nickname set "$NICKNAME" 2>&1 || true
fi
configure_peer_permission() {
    local t="$1" peers="$2" label="$3"
    if [ "$peers" = "*" ]; then
        log "Allowing all peers $label"
        nordvpn meshnet peer allow "$t" 2>&1 || true
    elif [ -n "$peers" ]; then
        log "Setting $label for: $peers"
        IFS=',' read -ra LIST <<< "$peers"
        for p in "${LIST[@]}"; do
            p=$(echo "$p" | xargs)
            [ -n "$p" ] && nordvpn meshnet peer allow "$t" "$p" 2>&1 || true
        done
    fi
}
configure_peer_permission "routing" "$ALLOW_ROUTING" "routing"
configure_peer_permission "local" "$ALLOW_LOCAL" "local network"
configure_peer_permission "fileshare" "$ALLOW_FILESHARE" "fileshare"
configure_peer_permission "remote" "$ALLOW_REMOTE" "remote access"
log "Verifying Meshnet..."
sleep 5
STATUS=$(nordvpn status 2>&1 || true)
log "Status: $STATUS"
log "Meshnet setup complete"
exit 0