# rootfs/etc/services.d/nordvpnd/run
#!/usr/bin/with-contenv bash
set -euo pipefail

echo "=== NordVPN Meshnet Starting ==="

# NordVPN needs a stable machine ID; HA add-on containers do not always provide one.
if [ ! -s /data/machine-id ]; then
    od -An -N16 -tx1 /dev/urandom | tr -d ' \n' > /data/machine-id
fi
cp /data/machine-id /etc/machine-id
mkdir -p /var/lib/dbus
cp /data/machine-id /var/lib/dbus/machine-id

# A failed previous service start can leave nordvpnd running and make the next
# start fail with "daemon is already running".
if pgrep -x nordvpnd >/dev/null 2>&1; then
    echo "Stopping existing NordVPN daemon..."
    pkill -TERM -x nordvpnd || true
    sleep 2
    pkill -KILL -x nordvpnd || true
fi

# Start the NordVPN daemon
echo "Starting NordVPN daemon..."
nordvpnd &
DAEMON_PID=$!
trap 'kill "$DAEMON_PID" 2>/dev/null || true' EXIT INT TERM
echo "NordVPN daemon started (PID: $DAEMON_PID)"

# Wait for the daemon socket to be ready. Do not use `nordvpn status` here:
# on a fresh install NordVPN prompts for analytics consent, which blocks forever
# in a non-interactive add-on container.
for i in $(seq 1 30); do
    if [ -S /run/nordvpn/nordvpnd.sock ]; then
        echo "NordVPN daemon is ready"
        break
    fi
    sleep 1
done

# Pre-seed analytics consent so later `nordvpn` CLI calls do not block waiting
# for interactive input. "off" sends only the essential data required to work.
printf 'n\n' | nordvpn set analytics off 2>&1 || true

# Authenticate if token is provided — read from /data/options.json (bashio unavailable with init: false)
CONFIG_FILE="/data/options.json"
TOKEN=""
if [ -f "$CONFIG_FILE" ]; then
    TOKEN="$(jq -r '.token // ""' "$CONFIG_FILE" 2>/dev/null || true)"
fi
if [ -n "$TOKEN" ]; then
    echo "Authenticating with NordVPN..."
    nordvpn login --token "$TOKEN" 2>&1 || echo "Auth failed (maybe already authenticated)"
    echo "Enabling Meshnet..."
    nordvpn set meshnet on 2>&1 || true
    NICKNAME="$(jq -r '.nickname // ""' "$CONFIG_FILE" 2>/dev/null || true)"
    if [ -n "$NICKNAME" ]; then
        nordvpn nickname set "$NICKNAME" 2>&1 || true
    fi
    nordvpn status 2>&1 || true
fi

# Start Web UI
echo "Setting up Web UI..."
mkdir -p /usr/src
cat > /usr/src/index.html << 'HTML'
<!DOCTYPE html><html lang="en"><head><meta charset="UTF-8"><meta name="viewport" content="width=device-width,initial-scale=1.0">
<title>NordVPN Meshnet</title>
<style>:root{--bg:#1a1a2e;--card:#16213e;--accent:#0f3460;--text:#e4e4e4;--success:#4ecca3;--error:#e94560}
*{margin:0;padding:0;box-sizing:border-box}
body{font-family:-apple-system,BlinkMacSystemFont,'Segoe UI',Roboto,sans-serif;background:var(--bg);color:var(--text);padding:2rem}
h1{text-align:center;color:var(--success);margin-bottom:2rem}
.card{background:var(--card);border-radius:12px;padding:1.5rem;border:1px solid var(--accent)}
.card div{padding:0.5rem 0}
.label{color:#888}.connected{color:var(--success)}.disconnected{color:var(--error)}
</style></head><body>
<h1>🔒 NordVPN Meshnet</h1>
<div class="card">
<div><span class="label">Status: </span><span id="s">Loading...</span></div>
<div><span class="label">Nord Name: </span><span id="n">—</span></div>
<div><span class="label">Meshnet IP: </span><span id="i">—</span></div>
</div>
<script>
fetch('/api/status').then(r=>r.json()).then(d=>{
document.getElementById('s').textContent=d.meshnet_status;document.getElementById('s').className=d.connected?'connected':'disconnected'
document.getElementById('n').textContent=d.meshnet_name;document.getElementById('i').textContent=d.meshnet_ip
}).catch(e=>document.getElementById('s').textContent='Error')
</script></body></html>
HTML

# Status API — proper HTTP response via ncat
cat > /usr/src/status_api.sh << 'STATUS_API'
#!/usr/bin/env bash
set -euo pipefail

serve_status() {
    local status m i
    status=$(nordvpn status 2>/dev/null || true)

    printf 'HTTP/1.1 200 OK\r\nContent-Type: application/json\r\nConnection: close\r\n\r\n'

    if echo "$status" | grep -qi "Status: Connected"; then
        m=$(echo "$status" | sed -n 's/.*Meshnet name:\s*//p' | tr -d '\r\n')
        i=$(echo "$status" | sed -n 's/.*Meshnet IP:\s*//p' | tr -d '\r\n')
        jq -nc --arg name "${m:-unknown}" --arg ip "${i:-unknown}" \
            '{connected:true,meshnet_name:$name,meshnet_ip:$ip,meshnet_status:"connected"}'
    else
        jq -nc '{connected:false,meshnet_name:"unknown",meshnet_ip:"unknown",meshnet_status:"disconnected"}'
    fi
}
serve_status
STATUS_API
chmod +x /usr/src/status_api.sh

ncat -l -k -p 8080 --sh-exec /usr/src/status_api.sh 2>/dev/null &
echo "Status API server started on port 8080"

# Start nginx for ingress
mkdir -p /etc/nginx/conf.d
cat > /etc/nginx/nginx.conf << 'NGINX_MAIN'
worker_processes 1;
events { worker_connections 128; }
http {
    include /etc/nginx/conf.d/*.conf;
}
NGINX_MAIN
cat > /etc/nginx/conf.d/default.conf << 'NGINX_SITE'
server {
    listen 80;
    location /api/status {
        proxy_pass http://127.0.0.1:8080;
    }
    location / {
        root /usr/src;
        index index.html;
    }
}
NGINX_SITE
echo "Starting nginx..."
nginx -t 2>&1 && nginx -g 'daemon off;' 2>&1 || echo "nginx failed to start"

# Keep the process alive
wait $DAEMON_PID