# meshnet/rootfs/etc/s6-overlay/services.d/nordvpn-setup/run
#!/usr/bin/with-contenv bash
# s6 service: One-time Meshnet setup
# Runs on first boot: authenticates, enables Meshnet, configures peer permissions

set -euo pipefail

log() {
    echo "[meshnet-setup] $1"
}

# Wait for daemon to be ready
log "Waiting for NordVPN daemon..."
for i in $(seq 1 30); do
    if nordvpn status >/dev/null 2>&1; then
        log "Daemon is ready"
        break
    fi
    sleep 1
done

# Read configuration
TOKEN=$(bash...n
# Validate token
if [ -z "$TOKEN" ]; then
    log "ERROR: No access token configured. Please add your token in the add-on options."
    exit 1
fi

# Step 1: Authenticate with NordVPN using token (headless login)
log "Authenticating with NordVPN..."
if ! nordvpn login --token "$TOKEN" 2>&1; then
    log "ERROR: Authentication failed. Check your token."
    exit 1
fi
log "Authentication successful"

# Step 2: Enable Meshnet mode (not regular VPN)
log "Enabling Meshnet mode..."
if ! nordvpn set meshnet on 2>&1; then
    log "ERROR: Failed to enable Meshnet"
    exit 1
fi
log "Meshnet mode enabled"

# Step 3: Set custom nickname (optional)
if [ -n "$NICKNAME" ]; then
    log "Setting device nickname: $NICKNAME"
    nordvpn nickname set "$NICKNAME" 2>&1 || true
fi

# Step 4: Configure peer permissions
configure_peer_permission() {
    local permission_type="$1"
    local peers="$2"
    local permission_label="$3"

    if [ "$peers" = "*" ]; then
        log "Allowing all peers $permission_label"
        nordvpn meshnet peer allow "$permission_type" 2>&1 || true
    elif [ -n "$peers" ]; then
        log "Setting $permission_label for: $peers"
        IFS=',' read -ra PEER_LIST <<< "$peers"
        for peer in "${PEER_LIST[@]}"; do
            peer=$(echo "$peer" | xargs)  # trim whitespace
            if [ -n "$peer" ]; then
                nordvpn meshnet peer allow "$permission_type" "$peer" 2>&1 || true
            fi
        done
    else
        log "No $permission_label peers configured"
    fi
}

configure_peer_permission "routing" "$ALLOW_ROUTING" "routing"
configure_peer_permission "local" "$ALLOW_LOCAL" "local network"
configure_peer_permission "fileshare" "$ALLOW_FILESHARE" "fileshare"
configure_peer_permission "remote" "$ALLOW_REMOTE" "remote access"

# Step 5: Verify Meshnet status
log "Verifying Meshnet status..."
sleep 5  # Give Meshnet time to initialize

STATUS=$(nordvpn status 2>&1 || true)
log "Meshnet status:"
echo "$STATUS"

# Get and display Meshnet details
NORD_NAME=$(nordvpn status 2>&1 | grep -oP 'Meshnet name: \K.*' || echo "unknown")
MESHNET_IP=$(nordvpn status 2>&1 | grep -oP 'Meshnet IP: \K.*' || echo "unknown")

log "Nord name: $NORD_NAME"
log "Meshnet IP: $MESHNET_IP"

log "Meshnet setup complete"

# Exit — s6 doesn't need this service to keep running (token persists in daemon state)
exit 0
